# contrib/hydrascale.service
[Unit]
Description=Hydrascale - Multi-Tailnet Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/hydrascale serve --config /etc/hydrascale/config.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
WorkingDirectory=/var/lib/hydrascale

# Minimal sandboxing — iptables-legacy and network namespace management
# require full privilege on some kernels (e.g. Tegra/Jetson), so we avoid
# NoNewPrivileges and ProtectSystem=strict which break iptables and netns.
ProtectHome=yes

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hydrascale

[Install]
WantedBy=multi-user.target
